# packaging/systemd/veya-daemon.service
# Veya daemon — systemd user service.
#
# Install (development build):
#   dotnet publish src/Veya.Daemon -c Release -o ~/.local/lib/veya
#   mkdir -p ~/.config/systemd/user
#   cp packaging/systemd/veya-daemon.service ~/.config/systemd/user/
#   systemctl --user daemon-reload
#   systemctl --user enable --now veya-daemon
#
# Logs: journalctl --user -u veya-daemon -f
#
# See docs/running.md for details.

[Unit]
Description=Veya intelligence daemon
Documentation=https://github.com/s3ba-b/ubuntu-intelligence

[Service]
Type=notify
ExecStart=%h/.local/lib/veya/Veya.Daemon
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
